# PROJETOS/20241020-portable-systemd/rootfs/etc/systemd/system/teste-main@.service
[Unit]
After=network.target
Description=Socket activated test with systemd-portable

[Service]
ExecStart=/teste
# StandardError=journal
# StandardInput=socket
# StandardOutput=socket
