# PROJETOS/20241020-portable-systemd/rootfs/etc/systemd/system/teste-main.socket
[Unit]

[Socket]
ListenStream=42069
Accept=yes

[Install]
WantedBy=sockets.target multi-user.target
